# package/deb/etc/systemd/system/casanode-startup.service
[Unit]
Description=Casanode Startup Initialization
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/casanode/startup.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
